# neural-tools/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:mcp-server]
command=python3 -u /app/src/mcp/neural_server_stdio.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/mcp-server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/mcp-server.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
priority=10
environment=PROJECT_NAME="%(ENV_PROJECT_NAME)s",PYTHONUNBUFFERED="1"

[program:indexer]
command=python3 -u /app/src/graphrag/services/multi_project_indexer.py --workspace-root=/workspace
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/indexer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/indexer.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
priority=20
environment=PROJECT_NAME="%(ENV_PROJECT_NAME)s",PYTHONUNBUFFERED="1"
startsecs=10
stopwaitsecs=30

[group:neural]
programs=mcp-server,indexer
priority=999